# debian/php5.3-fpm.upstart
# php5.3-fpm - The PHP FastCGI Process Manager

description "The PHP-5.3 FastCGI Process Manager"
author "Ondřej Surý <ondrej@debian.org>"

start on runlevel [2345]
stop on runlevel [016]

reload signal USR2

pre-start exec /usr/lib/php/5.3/php5.3-fpm-checkconf

respawn
exec /usr/sbin/php5.3-fpm --nodaemonize --fpm-config /etc/php/5.3/fpm/php.ini
